# deploy/services/django/run
#!/bin/bash
exec uwsgi --ini /src/project/deploy/uwsgi.ini